# crontab
44 1 * * * crontab < ~/highrise/crontab
@daily cd ~/highrise && ./scraper.py
